# ios-garment-measure/server/garment-3d-laptop.service
[Unit]
Description=Garment Seed3D API on Laptop (Tailscale)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=john
WorkingDirectory=/home/john/workspaces/search_similar_style/ios-garment-measure/server
Environment=ARK_API_KEY=ark-***
Environment=PUBLIC_BASE_URL=https://api.openfire.cloud
Environment=ARK_3D_MODEL=doubao-seed3d-2-0-260328
Environment=ARK_3D_SUBDIVISION=low
Environment=ARK_3D_FILE_FORMAT=usdz
Environment=ARK_3D_POLL_INTERVAL=30
Environment=ARK_3D_TIMEOUT_SEC=1800
ExecStart=/home/john/workspaces/search_similar_style/.venv/bin/python -m uvicorn garment_ai_server:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
